# cf-ddns.timer
[Unit]
Description=Timer for cf-ddns.service

[Timer]
OnBootSec=300
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
